# runit/updater/run
#!/bin/sh

exec 2>&1
cd /u/apps/rasplex-update-server/current

export UPDATER_ENVIRONMENT=production 
exec chpst -u deploy bundle exec ruby updateserver.rb >> /var/log/update-server
